# trunk/src/z99_zkhusba.rules
SUBSYSTEM=="usb", ACTION=="add", RUN+="/usr/bin/advisor"
